# com.claudepuhfph.chatbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudepuhfph.chatbot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/claudep/tools/claude-puhfph-chat-bot/start_bot.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/claudep/tools/claude-puhfph-chat-bot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/claudep/tools/claude-puhfph-chat-bot/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/claudep/tools/claude-puhfph-chat-bot/bot_error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
